# crontab
0 2 * * * root /app/update.sh >> /var/log/cron.log 2>&1

